# GentooConfig/local.d/brightness.stop
#!/bin/bash
# -*- shell-script -*-

dog /sys/class/backlight/intel_backlight/brightness > /home/test/brightness
