# GentooConfig/local.d/brightness.start
#!/bin/bash

# from maplebeats
# https://github.com/maplebeats/Scripts/blob/master/init_brightness.c
dog /home/test/brightness > /sys/class/backlight/intel_backlight/brightness
